# install/lss-backup.service
[Unit]
Description=LSS Backup Server
After=network.target mysql.service
Requires=mysql.service

[Service]
Type=simple
User=lss-backup
Group=lss-backup
ExecStart=/usr/local/bin/lss-backup-server
WorkingDirectory=/etc/lss-backup
Environment=LSS_ENV=production
Environment=LSS_CONFIG=/etc/lss-backup/config.toml
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lss-backup

# Hardening
NoNewPrivileges=false
PrivateTmp=true
ReadWritePaths=/etc/lss-backup /var/log/lss-backup /var/lib/lss-backup

[Install]
WantedBy=multi-user.target
